# freebsd/rc.d/dnsmasq_multi
#!/bin/sh
#
# PROVIDE: dnsmasq_multi
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Multiple dnsmasq instances for FreeBSD
#
# Add to /etc/rc.conf:
#   dnsmasq_multi_enable="YES"
#   dnsmasq_multi_instances="1 2 3 4 5 6 7 8"
#
# Each instance needs a config file:
#   /usr/local/etc/dnsmasq/dnsmasq1.conf
#   /usr/local/etc/dnsmasq/dnsmasq2.conf
#   etc.
#
# Each config must have unique:
#   - listen-address or interface
#   - pid-file (e.g., /var/run/dnsmasq/dnsmasq1.pid)
#   - port (if all on same IP)

. /etc/rc.subr

name="dnsmasq_multi"
rcvar="${name}_enable"

load_rc_config $name

: ${dnsmasq_multi_enable:="NO"}
: ${dnsmasq_multi_instances:="1 2 3 4 5 6 7 8"}
: ${dnsmasq_multi_bin:="/usr/local/sbin/dnsmasq"}
: ${dnsmasq_multi_confdir:="/usr/local/etc/dnsmasq"}
: ${dnsmasq_multi_piddir:="/var/run/dnsmasq"}

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"

dnsmasq_multi_start()
{
    echo "Starting dnsmasq instances..."

    # Erstelle PID-Verzeichnis
    mkdir -p "${dnsmasq_multi_piddir}"

    for i in ${dnsmasq_multi_instances}; do
        conf="${dnsmasq_multi_confdir}/dnsmasq${i}.conf"
        pid="${dnsmasq_multi_piddir}/dnsmasq${i}.pid"

        if [ ! -f "$conf" ]; then
            echo "  Instance $i: Config $conf nicht gefunden, überspringe"
            continue
        fi

        if [ -f "$pid" ] && kill -0 $(cat "$pid") 2>/dev/null; then
            echo "  Instance $i: Läuft bereits (PID $(cat $pid))"
            continue
        fi

        ${dnsmasq_multi_bin} --conf-file="$conf" --pid-file="$pid"

        if [ $? -eq 0 ]; then
            echo "  Instance $i: Gestartet"
        else
            echo "  Instance $i: Fehler beim Start"
        fi
    done
}

dnsmasq_multi_stop()
{
    echo "Stopping dnsmasq instances..."

    for i in ${dnsmasq_multi_instances}; do
        pid="${dnsmasq_multi_piddir}/dnsmasq${i}.pid"

        if [ -f "$pid" ]; then
            if kill -0 $(cat "$pid") 2>/dev/null; then
                kill $(cat "$pid")
                echo "  Instance $i: Gestoppt"
            fi
            rm -f "$pid"
        fi
    done
}

dnsmasq_multi_status()
{
    echo "Status dnsmasq instances:"

    for i in ${dnsmasq_multi_instances}; do
        pid="${dnsmasq_multi_piddir}/dnsmasq${i}.pid"

        if [ -f "$pid" ] && kill -0 $(cat "$pid") 2>/dev/null; then
            echo "  Instance $i: Läuft (PID $(cat $pid))"
        else
            echo "  Instance $i: Gestoppt"
        fi
    done
}

dnsmasq_multi_restart()
{
    dnsmasq_multi_stop
    sleep 1
    dnsmasq_multi_start
}

run_rc_command "$1"
